# docker/php-fpm/crontab/root
# -rw-------    1 root     root
# do daily/weekly/monthly maintenance   run-parts /etc/periodic/weekly
# min   hour    day     month   weekday command
*/15    *       *       *       *       run-parts /etc/periodic/15min
0       *       *       *       *       run-parts /etc/periodic/hourly
0       2       *       *       *       run-parts /etc/periodic/daily
0       3       *       *       6       run-parts /etc/periodic/weekly
0       5       1       *       *       run-parts /etc/periodic/monthly

# Laravel task schedule
* * * * * cd /var/www/html && php artisan schedule:run >> /dev/null 2>&1
